# config/systemd/user/hyprpolkitagent.service.d/override.conf
[Unit]
After=
After=dbus.service
PartOf=

[Install]
WantedBy=default.target
